# crontab
59 22 * * * python /app/daily_target.py
